# examples/personal-community-sentiment-triage/scripts/autoheal/systemd/nemoclaw-hermes-watchdog.timer.in
# SPDX-FileCopyrightText: Copyright (c) 2026 NVIDIA CORPORATION & AFFILIATES. All rights reserved.
# SPDX-License-Identifier: Apache-2.0
[Unit]
Description=Run the NemoClaw Hermes watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=nemoclaw-hermes-watchdog.service

[Install]
WantedBy=timers.target
